# deploy/anyarchie-worker.service
[Unit]
Description=AnyArchie Background Worker (Reminders)
After=network.target postgresql.service

[Service]
Type=simple
User=defibeats
WorkingDirectory=/var/www/anyarchie
ExecStart=/var/www/anyarchie/.venv/bin/python worker.py
Restart=on-failure
RestartSec=5
EnvironmentFile=/var/www/anyarchie/.env

[Install]
WantedBy=multi-user.target
