# deploy/anyarchie-bot.service
[Unit]
Description=AnyArchie Telegram Bot
After=network.target postgresql.service

[Service]
Type=simple
User=defibeats
WorkingDirectory=/var/www/anyarchie
ExecStart=/var/www/anyarchie/.venv/bin/python -u -m bot.main
Restart=on-failure
RestartSec=10
EnvironmentFile=/var/www/anyarchie/.env
StandardOutput=journal
StandardError=journal
KillMode=mixed
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
